# thermo/thermo-webserver.service
[Unit]
Description=Thermo sensor web server

[Service]
ExecStart=/opt/src/thermo/start_webserver.sh

[Install]
WantedBy=multi-user.target
